# timesync/timesync.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
55 17 * * * root /bin/bash /usr/bin/timesync
